# connman/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/connmand
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/connmanctl
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/connman-vpn-dbus.conf
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/connman.conf
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/connman/
-rw-r--r--	root/root	usr/include/connman/agent.h
-rw-r--r--	root/root	usr/include/connman/device.h
-rw-r--r--	root/root	usr/include/connman/inet.h
-rw-r--r--	root/root	usr/include/connman/inotify.h
-rw-r--r--	root/root	usr/include/connman/ipaddress.h
-rw-r--r--	root/root	usr/include/connman/ipconfig.h
-rw-r--r--	root/root	usr/include/connman/log.h
-rw-r--r--	root/root	usr/include/connman/machine.h
-rw-r--r--	root/root	usr/include/connman/network.h
-rw-r--r--	root/root	usr/include/connman/notifier.h
-rw-r--r--	root/root	usr/include/connman/peer.h
-rw-r--r--	root/root	usr/include/connman/plugin.h
-rw-r--r--	root/root	usr/include/connman/provision.h
-rw-r--r--	root/root	usr/include/connman/resolver.h
-rw-r--r--	root/root	usr/include/connman/service.h
-rw-r--r--	root/root	usr/include/connman/session.h
-rw-r--r--	root/root	usr/include/connman/storage.h
-rw-r--r--	root/root	usr/include/connman/version.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/connman/
drwxr-xr-x	root/root	usr/lib/connman/plugins-vpn/
-rwxr-xr-x	root/root	usr/lib/connman/plugins-vpn/openvpn.la
-rwxr-xr-x	root/root	usr/lib/connman/plugins-vpn/openvpn.so
drwxr-xr-x	root/root	usr/lib/connman/scripts/
-rwxr-xr-x	root/root	usr/lib/connman/scripts/openvpn-script
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/connman.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/connman-vpnd
-rwxr-xr-x	root/root	usr/sbin/connmand
-rwxr-xr-x	root/root	usr/sbin/connmand-wait-online
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/net.connman.vpn.service
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/connmanctl.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/connman-service.config.5.gz
-rw-r--r--	root/root	usr/share/man/man5/connman-vpn-provider.config.5.gz
-rw-r--r--	root/root	usr/share/man/man5/connman-vpn.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/connman.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/connman-vpn.8.gz
-rw-r--r--	root/root	usr/share/man/man8/connman.8.gz
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/net.connman.policy
-rw-r--r--	root/root	usr/share/polkit-1/actions/net.connman.vpn.policy
